# com.bryanrowland.brytools-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bryanrowland.brytools-watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Volumes/RowMedia/CODE/brytools/watchdog.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/bryanrowland/.brytools-watchdog/watchdog-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bryanrowland/.brytools-watchdog/watchdog-stderr.log</string>
</dict>
</plist>
